# bacula/platforms/systemd/bacula-sd.service.in
# Systemd Bacula service file
#
# Copyright (C) 2000-2022 Kern Sibbald
# License: BSD 2-Clause; see file LICENSE-FOSS
#
# /lib/systemd/system/bacula-sd.service
#
# Description:
#    Used to start the bacula storage daemon service (bacula-sd)
#    enable : systemctl enable bacula-sd
#    start : systemctl start bacula-sd
#
#

# from http://www.freedesktop.org/software/systemd/man/systemd.unit.html
[Unit]
Description=Bacula Storage Daemon service
Requires=network.target
After=network.target
RequiresMountsFor=@working_dir@ @sysconfdir@ @sbindir@

# from http://www.freedesktop.org/software/systemd/man/systemd.service.html
[Service]
Type=exec
User=@sd_user@
Group=@sd_group@
SupplementaryGroups=@sd_user@
Environment="LD_LIBRARY_PATH=@libdir@"
# comment out to get traces working at startup using -T and -d options
#WorkingDirectory=@working_dir@
ExecStart=@sbindir@/bacula-sd -fP -c @sysconfdir@/bacula-sd.conf
StandardError=syslog
TimeoutStopSec=3min
LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target
